# roles/inkjet-maintenance/templates/etc/systemd/system/inkjet-maintenance.timer
[Unit]
Description=Print something every few days

[Timer]
OnUnitActiveSec=3d

[Install]
WantedBy=timers.target